# main/docker/monitoringstack/services/scrutiny.service
[Unit]
Description=Scrutiny Collector
Requires=scrutiny.timer

[Service]
Type=simple
ExecStart=/node/scripts/scrutiny.sh
User=root
Group=root
